# systemd/linux-voice-assistant.service
[Unit]
Description=Linux Voice Assistant (Samantha)
After=network.target sound.target voice-orb-bridge.service
Wants=voice-orb-bridge.service

[Service]
User=satellite
WorkingDirectory=/home/satellite/linux-voice-assistant
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PULSE_SERVER=unix:/run/user/1000/pulse/native
ExecStart=/home/satellite/linux-voice-assistant/.venv/bin/python3 \
  -m linux_voice_assistant \
  --name Samantha \
  --wake-word-dir /home/satellite/voice_orb/models \
  --wake-word-dir /home/satellite/linux-voice-assistant/wakewords \
  --wake-model samantha \
  --audio-input-device "Built-in Audio Pro" \
  --audio-output-device "pipewire/alsa_output.platform-soc_sound.pro-output-0" \
  --debug
Restart=always
RestartSec=8

[Install]
WantedBy=multi-user.target
