# systemd/linux-voice-assistant.service
[Unit]
Description=Linux Voice Assistant (Samantha)
Documentation=https://github.com/OHF-Voice/linux-voice-assistant
After=network.target sound.target voice-orb-bridge.service
Wants=voice-orb-bridge.service

[Service]
User=pi
WorkingDirectory=/home/pi/linux-voice-assistant

ExecStart=/home/pi/linux-voice-assistant/.venv/bin/python3 \
  -m linux_voice_assistant \
  --name Samantha \
  --wake-word-dir /home/pi/voice-orb/models \
  --wake-model samantha \
  --audio-input-device "YOUR_MIC_DEVICE_NAME" \
  --audio-output-device "YOUR_SPEAKER_DEVICE_NAME"

# Pipe output to the log file the bridge watches
StandardOutput=append:/tmp/lva.log
StandardError=append:/tmp/lva.log

Restart=always
RestartSec=5
# Give audio subsystem time to start
ExecStartPre=/bin/sleep 3

[Install]
WantedBy=multi-user.target
