# systemd/voice-orb-bridge.service
[Unit]
Description=Voice Orb Bridge (Samantha)
Documentation=https://github.com/OHF-Voice/linux-voice-assistant
After=network.target sound.target

[Service]
User=satellite
WorkingDirectory=/home/satellite/voice_orb

ExecStart=/usr/bin/python3 /home/satellite/voice_orb/voice_orb_bridge.py

Restart=always
RestartSec=5
KillMode=control-group

[Install]
WantedBy=multi-user.target
